# scraper/crontab
*/5 * * * * /usr/local/bin/python /app/ingest_sessions.py
*/5 * * * * /usr/local/bin/python /app/ingest_usage.py
