# scripts/com.claude-code-analytics.up.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-code-analytics.up</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>for i in $(seq 1 60); do /opt/homebrew/bin/docker info >/dev/null 2>&1 && break; sleep 5; done; /Users/alejandrosanchez/Workspace/claude-code-analytics/scripts/up.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/claude-code-analytics-up.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-code-analytics-up.log</string>
</dict>
</plist>
